# addons/iot_box_image/overwrite_after_init/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

exec 1>/var/log/rc.local.log 2>&1  # send stdout and stderr from rc.local to a log file
set -x                             # display commands before execution

# Ensure Wi-Fi radio is enabled
nmcli radio wifi on

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

mkdir -p /var/run/betopiaerp
chown betopiaerp:betopiaerp /var/run/betopiaerp

# For compatibility between new IoT Box images and old BetopiaERP code on it (it will not start BetopiaERP on boot without this)
start_wifi=/home/pi/betopiaerp/addons/point_of_sale/tools/posbox/configuration/wireless_ap.sh
if [ -f $start_wifi ]; then
  $start_wifi &
fi

# Update current branch
cd /home/pi/betopiaerp
localbranch=$(sudo -u betopiaerp git symbolic-ref -q --short HEAD)
localremote=$(sudo -u betopiaerp git config branch.$localbranch.remote)

sudo -u betopiaerp rm -f /home/pi/betopiaerp/.git/shallow.lock
sudo -u betopiaerp git remote set-url "${localremote}" "https://github.com/betopiaerp/betopiaerp.git"

sudo -u betopiaerp GIT_SSL_NO_VERIFY=1 git fetch "${localremote}" "${localbranch}" --depth=1
sudo -u betopiaerp git reset --hard FETCH_HEAD

exit 0
